# bin/throttlestop.service
[Unit]
Description=ThrottleStop for Linux.

[Service]
Type=simple
User=root
ExecStart=/usr/bin/throttlestop

[Install]
WantedBy=multi-user.target
